# dot_files/emacs/dot_getmail.crontab
# m h  dom mon dow   command
*/30 * * * * $HOME/.getmail/getmail.sh

